# main/postmarketos-ui-retroarch/postmarketos-ui-retroarch.post-install
#!/bin/sh

rc-update add tinydm
rc-update add bluetooth default

tinydm-set-session -s /usr/share/wayland-sessions/retroarch.desktop
