# temp/hexagonrpcd/hexagonrpcd-adsp.initd
#!/sbin/openrc-run

supervisor=supervise-daemon
name="Hexagonrpcd ADSP"
description="Daemon to support Qualcomm Hexagon ADSP virtual filesystem"

command=/usr/sbin/hexagonrpcd
command_args="-f /dev/fastrpc-adsp"
