# main/nodm/nodm.openrc
#!/sbin/openrc-run

description="Automatic display manager"
command="/usr/sbin/nodm"
command_background="true"
pidfile="/var/run/nodm.pid"
